# deploy/supervisor.conf
[program:cineweave-api]
command=/var/www/cineweave/venv/bin/uvicorn server.app:app --host 127.0.0.1 --port 8000 --workers 2
directory=/var/www/cineweave
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/cineweave/error.log
stdout_logfile=/var/log/cineweave/access.log
environment=PATH="/var/www/cineweave/venv/bin",NODE_ENV="production",APP_ENV="production"
